# bin/randomWallpaper.cron
*/15 * * * * DESKTOP_SESSION=i3 DISPLAY=:0.0 $HOME/dotfiles/bin/randomWallpaper.sh